# com.omsk-bus-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.omsk-bus-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gleb/tools/omsk-bus-bot/venv/bin/python</string>
        <string>-m</string>
        <string>omsk_bus_bot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gleb/tools/omsk-bus-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/gleb/tools/omsk-bus-bot/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gleb/tools/omsk-bus-bot/logs/bot.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
